# systemd/takoyaki-predictor.service
[Unit]
Description=Takoyaki Predictor Dashboard (real detections)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=nishikohh
WorkingDirectory=/home/nishikohh/peopleflow_app
ExecStart=/home/nishikohh/peopleflow_app/scripts/run_predictor_real.sh
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
